# templates/vsts-agent.service.j2
[Unit]
Description=Azure Pipelines Agent
After=network.target

[Service]
ExecStart={{vsts_agentfolder}}/bin/runsvc.sh
User={{vsts_agent_user}}
WorkingDirectory={{vsts_agentfolder}}
KillMode=process
KillSignal=SIGTERM
TimeoutStopSec=5min

[Install]
WantedBy=multi-user.target